# baremetal_minio/restart_avahi.timer
[Unit]
Description="Run script after 3 minutes of boot"

[Timer]
OnBootSec=3min

[Install]
WantedBy=timers.target
