# stress_test/stress_cpu.service
[Unit]
Description=CPU_Stress_Test
[Service]
ExecStart=/usr/bin/stress --cpu 6
Restart=always
[Install]
WantedBy=multi-user.target
