# local_minio/minio_lookup.service
[Unit]
Description=MinioLookup
[Service]
ExecStart=/usr/bin/service_lookup.py
Restart=always
[Install]
WantedBy=multi-user.target
